# environment/supervisor.conf
[program:tomorrow]
command=nodejs app.js
directory=/home/melfa/tomorrow
user=melfa
stdout_logfile=/home/melfa/tomorrow/runtime/supervisor.log
